# services/php/finish
#!/bin/sh

kill `pidof php-fpm5`
[ -f /run/php/silo.sock ] && rm /run/php/silo.sock
