# services/syslog-ng/run
#!/bin/sh
exec 2>&1
source /etc/envvars

exec /usr/sbin/syslog-ng --pidfile /run/syslog-ng.pid --foreground
